# etc/hebcal-api.service
[Unit]
Description=Hebcal API server (Hebrew Date Converter + Zmanim REST APIs)
Documentation=https://github.com/hebcal/hebcal-api-go
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
ExecStart=/usr/local/bin/hebcal-api \
  -port 8082 \
  -logfile /var/log/hebcal/api.log \
  -zips-db /var/www/zips.sqlite3 \
  -geonames-db /var/www/geonames.sqlite3
# SIGUSR1 (or SIGHUP) reopens the access log after rotation
ExecReload=/bin/kill -USR1 $MAINPID
Restart=on-failure
RestartSec=2

# systemd creates /var/log/hebcal owned by the service user
LogsDirectory=hebcal
LogsDirectoryMode=0755

# hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
