# etc/hebcal-converter.service
[Unit]
Description=Hebcal Hebrew Date Converter REST API
Documentation=https://github.com/hebcal/converter-web
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
ExecStart=/usr/local/bin/hebcal-converter -port 8082 -logfile /var/log/hebcal/converter.log
# SIGUSR1 (or SIGHUP) reopens the access log after rotation
ExecReload=/bin/kill -USR1 $MAINPID
Restart=on-failure
RestartSec=2

# systemd creates /var/log/hebcal owned by the service user
LogsDirectory=hebcal
LogsDirectoryMode=0755

# hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
